# systemd/dreamconnect-daemon.service
[Unit]
Description=dreamconnect runtime daemon (Wayland capture/input bridge for ScreenConnect)
Documentation=https://github.com/ghostpsalm/DreamConnect
# Needs the GNOME/Wayland session (Mutter D-Bus, PipeWire) to be up.
After=graphical-session.target
PartOf=graphical-session.target
# Retry forever: never let a fast crash-loop trip the default start rate-limit
# (5 starts / 10s) and leave the bridge permanently stopped on an unattended box.
StartLimitIntervalSec=0

[Service]
Type=simple
# @INSTALL_DIR@ and @MONITOR@ are substituted by install.sh.
ExecStart=/usr/bin/python3 @INSTALL_DIR@/runtime/dreamconnect_daemon.py --monitor @MONITOR@
Restart=always
RestartSec=2
# Owner-only files: the root SC JVM reads the shm/socket via DAC override, so
# other local users can't scrape the screen or inject input.
UMask=0077

[Install]
WantedBy=graphical-session.target
